# package/base-files/files/etc/crontabs/root

10 6 * * * sleep 5 && touch /etc/banner && reboot
0 */1 * * * curl -A " https://baidu.com
